# net-audit-export.service
[Unit]
Description=Export net_connect audit log entries to user-readable CSV

[Service]
Type=oneshot
ExecStart=__INSTALL_DIR__/export_net_audit.sh
